# net/Murillo-pkg-captiveAlias/files/etc/cron.d/captivealias.cron
*	*	*	*	*	root	/usr/bin/nice -n20 /usr/local/bin/php -r 'require_once("captivealias.inc"); captivealias_cron_sync();'
