# net/Murillo-pkg-captiveAlias/pkg-plist
pkg/captivealias.xml
pkg/captivealias.inc
/etc/cron.d/captivealias.cron
%%DATADIR%%/info.xml
